# 普罗米修斯+node_exporter自动化/node_exporter.service
[Unit]
Description=node exporter

[Service]
ExecStart=/usr/local/bin/node_exporter/node_exporter
 
[Install]
WantedBy=multi-user.target
