# 普罗米修斯+node_exporter自动化/prometheus.service
[Unit]
Description=prometheus service

[Service]
ExecStart=/usr/local/bin/prometheus/prometheus --config.file=/usr/local/bin/prometheus/prometheus.yml 
 
[Install]
WantedBy=multi-user.target
